# deploy/docker/crontab
SHELL=/bin/sh
PATH=/app/.venv/bin:/usr/local/bin:/usr/bin:/bin

# 所有任务由 writer lock 串行保护；失败退出码保留给容器日志和监控。
30 14 * * 1-5 root cd /app && quant data daily --trade-date today >> /app/logs/daily.log 2>&1
20 15 * * 1-5 root cd /app && quant data daily --trade-date today >> /app/logs/daily.log 2>&1
30 18 * * 1-5 root cd /app && quant data daily --trade-date today >> /app/logs/daily.log 2>&1
30 19 * * 1-5 root cd /app && quant data backfill --start today --end latest >> /app/logs/backfill.log 2>&1
30 22 * * 1-5 root cd /app && quant data reconcile --trade-date today >> /app/logs/reconcile.log 2>&1
30 7 * * 1-5 root cd /app && quant data publish --trade-date today >> /app/logs/publish.log 2>&1
0 23 * * 5 root cd /app && quant data snapshot --profile dev >> /app/logs/snapshot.log 2>&1
30 23 * * 0 root cd /app && quant data snapshot --profile full >> /app/logs/snapshot.log 2>&1
0 2 1 * * root cd /app && quant data compact >> /app/logs/compact.log 2>&1
